# ansible/roles/smtp-server/templates/smtp-server.service.j2
[Unit]
Description=SMTP Server
Documentation=https://github.com/your-org/smtp-server
After=network.target

[Service]
Type=simple
User={{ deploy_user }}
Group={{ deploy_group }}
EnvironmentFile={{ config_dir }}/smtp.env
ExecStart={{ install_dir }}/smtp-server
ExecReload=/bin/kill -HUP $MAINPID
Restart={{ 'on-failure' if restart_on_failure else 'no' }}
RestartSec={{ restart_delay_sec }}

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ data_dir }} {{ queue_dir }} {{ mail_dir }} {{ log_dir }}

# Limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=smtp-server

[Install]
WantedBy=multi-user.target
